# files/etc/uci-defaults/99_gateway_watchdog
#!/bin/sh
# 99_gateway_watchdog — uci-defaults first-boot initialisation
# Runs once automatically via /etc/init.d/boot, then is deleted.
# When run manually by install.sh, install.sh deletes it afterward.

# -------------------------------------------------------
# Initialise UCI config if not already present
# -------------------------------------------------------
uci -q get gateway_watchdog.settings || {
    uci set gateway_watchdog.settings=settings
    uci set gateway_watchdog.settings.enabled=1
    uci set gateway_watchdog.settings.interface=wan
    uci set gateway_watchdog.settings.delay=10
    uci set gateway_watchdog.settings.max_failures=5
    uci set gateway_watchdog.settings.cooldown=300
    uci set gateway_watchdog.settings.recovery_mode=full
    uci set gateway_watchdog.settings.recovery_verify_targets="8.8.8.8,1.1.1.1"
    uci set gateway_watchdog.settings.log_to_console=0
    uci commit gateway_watchdog
}

# -------------------------------------------------------
# Preserve user config and runtime scripts across sysupgrade.
#
# Only config and user-data files belong in sysupgrade.conf.
# Package-managed files (binaries, LuCI assets, rpcd scripts) must NOT
# be listed here — sysupgrade would restore the pre-upgrade version over
# the freshly installed one, breaking the service after every firmware update.
# -------------------------------------------------------
add_to_sysupgrade() {
    grep -q "^$1$" /etc/sysupgrade.conf 2>/dev/null || echo "$1" >> /etc/sysupgrade.conf
}

# User config — must survive sysupgrade
add_to_sysupgrade "/etc/config/gateway_watchdog"

# Runtime data — optional but useful to preserve across upgrades
add_to_sysupgrade "/tmp/gateway_watchdog.history"
add_to_sysupgrade "/tmp/gateway_watchdog_stats"

# NOT preserved (package-managed, reinstalled by install.sh after upgrade):
#   /etc/init.d/gateway_watchdog
#   /etc/hotplug.d/iface/30-gateway_watchdog
#   /usr/bin/gateway_watchdog.sh
#   /usr/libexec/rpcd/gateway-watchdog-status
#   /usr/share/rpcd/acl.d/luci-app-gateway-watchdog.json
#   /usr/share/luci/menu.d/luci-app-gateway-watchdog.json
#   /www/luci-static/resources/view/gateway-watchdog/

# -------------------------------------------------------
# Enable the service so the S95 symlink is created.
#
# Do NOT call start here. uci-defaults runs early in the boot
# sequence before networking is initialised. The S95gateway_watchdog
# symlink will start the service at the correct time (priority 95,
# after network at 20). When run by install.sh directly, install.sh
# prints a reminder to start the service manually if needed.
# -------------------------------------------------------
/etc/init.d/gateway_watchdog enable

exit 0
